# example.ini
[program:example]
command=/bin/bash /mnt/noexit.sh
directory=/mnt
user=root
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/data/logs/example.std.log
loglevel=info

